# ops/systemd/shkandal-remote-ml.service
[Unit]
Description=Run one Shkandal remote ML worker pass
Requires=shkandal-remote-db-tunnel.service
After=shkandal-remote-db-tunnel.service

[Service]
Type=oneshot
WorkingDirectory=/opt/shkandal
ExecStart=/usr/bin/env bash ops/run-remote-worker worker-ml
TimeoutStartSec=30min
