# ops/systemd/shkandal-ml-worker.timer
[Unit]
Description=Run Shkandal ML worker five minutes after the previous pass finishes

[Timer]
OnBootSec=2min
OnUnitInactiveSec=5min
Persistent=true
Unit=shkandal-ml-worker.service

[Install]
WantedBy=timers.target
