# ops/systemd/shkandal-ingestion.service
[Unit]
Description=Run one Shkandal ingestion pass
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/shkandal
ExecStart=/usr/bin/env bash ops/run-scheduled-worker worker-ingestion
TimeoutStartSec=45min
